# configurations/etc/systemd/system/docker.service.d/http-proxy.conf
[Service]
Environment="HTTP_PROXY=socks5://127.0.0.1:2080"
Environment="HTTPS_PROXY=socks5://127.0.0.1:2080"
Environment="NO_PROXY=.ficld.ir,localhost,127.0.0.1"